# service/app/dieklingel-core.service.in
[Unit]
Description=dieKlingel Core

[Service]
Type=simple
Restart=always
WorkingDirectory=/etc/dieklingel-core
ExecStart=@EXE_PATH@

[Install]
WantedBy=mulit-user.target
